# green-proxy.service
[Unit]
Description=green-proxy: Reverse Proxy Server

[Service]
ExecStart=~/green-proxy
Restart=always

[Install]
WantedBy=multi-user.target
