# green-proxy-linux.service
[Unit]
Description=Reverse Proxy Server

[Service]
ExecStart=/etc/green-proxy/proxy-runner
Restart=always

[Install]
WantedBy=multi-user.target
